# etc/systemd/system/kernel-boot-update.path
Description=Copies Kernel and RamDisk to the boot partition

[Path]
PathChanged=/boot/initramfs-linux-t2-fallback.img

[Install]
WantedBy=multi-user.target

